# packaging/systemd/modelmux.service
[Unit]
Description=ModelMux - OpenAI-compatible proxy for Vertex AI
Documentation=https://github.com/yarenty/modelmux
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/modelmux
Restart=on-failure
RestartSec=5

# Configuration: reads from /etc/modelmux/config.toml (system) or
# ~/.config/modelmux/config.toml (user config)
# The config.toml file should contain:
#   [server] - port, log_level, retries
#   [auth] - service_account_json or service_account_file
#   [vertex] - project, region, location, publisher, model
#   [streaming] - mode, buffer_size, chunk_timeout_ms

# Security hardening
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/etc/modelmux

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=modelmux

# Performance
LimitNOFILE=65536

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
